# .config/autostart/Enpass.desktop
[Desktop Entry]
Type=Application
Name=Enpass
Exec=/opt/enpass/Enpass -minimize
X-GNOME-Autostart-enabled=true
