# .config/autostart/signal.desktop
[Desktop Entry]
Type=Application
Name=Signal
Exec=signal-desktop --password-store=basic --start-in-tray
X-GNOME-Autostart-enabled=true
